# scripts/vibe-anywhere.service
[Unit]
Description=Vibe Anywhere - AI Coding Session Manager
Documentation=https://github.com/kobozo/vibe-anywhere
After=network.target postgresql.service

[Service]
Type=simple
User=vibeanywhere
Group=vibeanywhere
WorkingDirectory=/opt/vibe-anywhere
EnvironmentFile=/opt/vibe-anywhere/.env
ExecStart=/usr/bin/node server.js
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/opt/vibe-anywhere /home/vibeanywhere
PrivateTmp=true

[Install]
WantedBy=multi-user.target
